# br2-external/board/f1c200s-devboard/rootfs-overlay/etc/inittab
::sysinit:/etc/init.d/rcS
ttyS0::respawn:-/bin/sh
ttyGS0::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
